# srcs/requirements/wordpress/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm83]
command=php-fpm83 -F
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.err.log
stdout_logfile=/var/log/php-fpm.out.log